# etc/smspipeline.conf
[program:smspipeline]
command=python /opt/smspipeline/web/manage.py runserver 0.0.0.0:80 --settings=settings-aws
logfile=/var/log/supervisor/smspipeline.log
log_stdout=true
log_stderr=true
